# xap-usb-gadget/systemd/xap-usb-gadget.service.in
[Unit]
Description=Set up USB gadget on wlan0

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'export MAC=`/bin/cat /sys/class/net/wlan0/address`; /sbin/modprobe g_ether host_addr=$MAC'
ExecStartPost=/bin/sh -c 'sleep 1'
ExecStop=/sbin/modprobe -r g_ether
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
