# xap-usb-gadget/systemd/xap-usb-gadget.service.in
[Unit]
Description=Set up USB gadget on wlan0
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=xap-usb-gadget create
ExecStartPost=/bin/sh -c 'sleep 1'
ExecStop=xap-usb-gadget destroy
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
