# xbridge/systemd/xbridge.service.in
[Unit]
Description=xbridge service (pass-through mode)
After=xap-usb-gadget.service

[Service]
Type=simple
ExecStartPre=/sbin/ip link set dev usb0 up
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/xbridge
ExecStopPost=/sbin/ip link set dev usb0 down
ProtectSystem=true
ProtectHome=true
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
