# xbridge/systemd/xbridges.service.in
[Unit]
Description=xbridge service (secure mode)
BindsTo=enftun-data@enf0.service
After=enftun-data@enf0.service xap-usb-gadget.service

[Service]
Type=simple
ExecStartPre=/sbin/ip link set dev usb0 up
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/xbridge -2 -n enf0
ExecStopPost=/sbin/ip link set dev usb0 down
ProtectSystem=true
ProtectHome=true
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

