# xbridge/systemd/xbridges-setup.service.in
[Unit]
Description=setup for xbridge service (secure mode)
Requires=connman.service
After=enftun@enf0.service
BindsTo=xbridges.service

[Service]
Type=oneshot
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/xbridge-setup start enf0
ExecStop=@CMAKE_INSTALL_FULL_BINDIR@/xbridge-setup stop enf0
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
